# addon/rootfs/etc/services.d/luftujha/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

bashio::log.info "Starting the Luftuj backend service..."

PORT=$(bashio::config 'web_port')
if [ -z "${PORT}" ] || [ "${PORT}" = "null" ]; then
  PORT=8099
fi

cd /usr/src/app
PORT="${PORT}" exec bun run src/server.ts
