# addon/rootfs/etc/services.d/luftujha/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

bashio::log.info "Starting the Luftuj backend service..."

PORT=$(bashio::config 'web_port')
if [ -z "${PORT}" ] || [ "${PORT}" = "null" ]; then
  PORT=8099
fi

function handle_error() {
    bashio::log.error "An error occurred while starting the Luftuj backend service."
    bashio::log.error "Error on line $1: Command '$2' exited with status $3"
    exit 1
}
trap 'handle_error $LINENO "$BASH_COMMAND" $?' ERR

cd /usr/src/app || bashio::exit.nok "Could not change directory to /usr/src/app"

bashio::log.info "Starting server on port ${PORT}..."
PORT="${PORT}" exec bun run src/server.ts || bashio::exit.nok "Failed to start the server"
